# app/backend/crontab
*/5 * * * * cd /usr/src/app && python3 data_pipeline.py >> /var/log/cron.log 2>&1