# app/backend/crontab
*/10 * * * * cd /usr/src/app && python3 data_pipeline.py >> /var/log/cron.log 2>&1
0 3 * * 1 docker builder prune -a -f
